# etc/upstart/wh-postgis.conf
description "Whathood Containers"
author "Jim Smiley"
start on filesystem and started docker
stop on runlevel [!2345]
respawn
script
  /usr/bin/docker start -a wh-postgis
end script
